# packages/operbot/operbot-106-py3-none-any.whl/operbot-106.data/data/share/operbot/operbot.service
# This file is placed in the Public Domain


[Unit]
Description=OPERBOT - write your own commands
After=multi-user.target

[Service]
Type=forking
KillMode=mixed
PrivateTmp=true
DynamicUser=True
StateDirectory=operbot
LogsDirectory=operbot
CacheDirectory=operbot
ExecStart=/usr/local/bin/operbotd
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
