# packages/rssbot/rssbot-123.tar.gz/rssbot-123/files/rssbot.service
[Unit]
Description=rssbot feeding rss into your irc channel
After=multi-user.target

[Service]
Type=forking
KillMode=mixed
PrivateTmp=true
DynamicUser=True
StateDirectory=rssbot
LogsDirectory=rssbot
CacheDirectory=rssbot
ExecStart=/usr/local/bin/rssbot
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
